# script/99-remove-usb.rules
ACTION=="remove",KERNEL=="sd[a-z][0-9]",RUN+="/opt/ros/indigo/bin/unmountusb.sh"
